# root.cron
32 20	*	*	*	/usr/local/bin/backup-mongodb
