# root.cron
* *	*	*	*	/usr/local/bin/backup-mongodb
